# services/aqch-start.service
[Unit]
Description=XPower AQCH miner (re)starter

[Service]
ExecStart=/usr/bin/systemctl restart --user aqch-miner@%i.service
Type=oneshot
