# services/loki-start@.service
[Unit]
Description=XPower LOKI miner (re)starter

[Service]
ExecStart=/usr/bin/systemctl restart --user loki-miner@%i.service
Type=oneshot
