# services/odin-start@.service
[Unit]
Description=XPower ODIN miner (re)starter

[Service]
ExecStart=/usr/bin/systemctl restart --user odin-miner@%i.service
Type=oneshot
