# services/xpow-start@.service
[Unit]
Description=XPower miner (re)starter

[Service]
ExecStart=/usr/bin/systemctl restart --user xpow-miner@%i.service
Type=oneshot
